# systemd/linuxhello-inference.service
[Unit]
Description=LinuxHello Inference Service (AI Backend)
After=network.target

[Service]
Type=simple
User=linuxhello
Group=linuxhello
SupplementaryGroups=video
WorkingDirectory=/usr/share/linuxhello/python-service
ExecStartPre=/usr/libexec/linuxhello/sync-python-venv.sh /usr/share/linuxhello/python-service
ExecStart=/usr/share/linuxhello/python-service/venv/bin/python3 inference_service.py --host 127.0.0.1 --port 50051
# Clean shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=10
Restart=on-failure
RestartSec=5s
UMask=0077

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadWritePaths=/usr/share/linuxhello /var/lib/linuxhello

[Install]
WantedBy=multi-user.target